# docker/supervisord.conf
; agent-hub bundle image: supervisord config (= issue #95、 Dockerfile.bundle 用)
;
; server (Node MCP) + scheduler (Python cron daemon) を同 container で並走させる
; ための supervisor program 定義。 supervisord を PID 1 (= `-n` foreground mode)
; として起動、 child processes の lifecycle 管理 + log forwarding を行う。
;
; 設計:
; - 両 process とも `autorestart=true` で transient failure に対する resilience
; - server (= MCP) は scheduler の親 (= scheduler が server に MCP 接続する peer)
;   のため、 `priority=1` で先に起動、 scheduler は `priority=2` で後追い
; - stdout / stderr は container log (= `docker logs`) に forward (= 12-factor
;   の log handling と整合、 container 内で file に書かない)

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid
; non-root への drop は base image (= node:20-bookworm-slim) の root ユーザー前提
; のままにする。 将来 security hardening 余地として別 PR で USER 切り替え検討。

[program:agent-hub-server]
; tsx で TS を直接実行 (= 既存 Dockerfile と同等の運用)。 production deployment では
; pre-compile (= tsc → node dist/) への切り替え余地あるが、 本 PR scope 外。
command=npm run mcp:start
directory=/app
priority=1
autostart=true
autorestart=true
startsecs=5
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production"

[program:agent-hub-scheduler]
; venv 経由で scheduler.py を起動 (= Dockerfile で /app/.venv/bin を PATH 通済)。
; AGENT_HUB_URL は ENV で渡し (= Dockerfile で localhost:3000 を default)、
; supervisord が server を先に起動するので scheduler が register / subscribe する
; タイミングで server が ready 想定 (= 万一 race ったら scheduler 内 retry が effective)。
command=/app/.venv/bin/python /app/packages/scheduler/scheduler.py
directory=/app/packages/scheduler
priority=2
autostart=true
autorestart=true
startsecs=10
stopsignal=TERM
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
; AGENT_HUB_URL は Dockerfile で ENV 済、 GITHUB_PAT / AGENT_HUB_USER /
; AGENT_HUB_TENANT は `docker run -e ...` で外部から渡す前提。
